# lib/systemd/system/ezekiel_updatedb.timer
[Unit]
Description=Run probe_updatedb every few seconds and on boot

[Timer]
OnBootSec=10s
OnUnitActiveSec=5s

[Install]
WantedBy=timers.target
