# lib/systemd/system/probe_updatedb.timer
[Unit]
Description=Run probe_updatedb every 3s and on boot

[Timer]
OnBootSec=10s
OnUnitActiveSec=3s 

[Install]
WantedBy=timers.target
